# Object_config/Storage_NFS/var-lib-glance-images.automount
[Unit]
Description=NFS Storage Automount
Wants=network.target rpcbind.service
After=network.service rpcbind.service
[Automount]
Where=/var/lib/glance/images/
[Install]
WantedBy=multi-user.target
