# Object_config/Storage_NFS/Nova_Api/var-lib-nova.automount
[Unit]
Description=NFS Storage Automount
Wants=network.target rpcbind.service
After=network.service rpcbind.service
[Automount]
Where=/var/lib/nova
[Install]
WantedBy=multi-user.target
